# private_dot_config/systemd/user/syncthing.service
[Unit]
Description=Syncthing

[Service]
ExecStart=/usr/bin/syncthing --no-browser
